# telegram-music-bot.service
[Unit]
Description=Telegram Music Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/telegram-music-bot
Environment=PATH=/home/ubuntu/telegram-music-bot/venv/bin
ExecStart=/home/ubuntu/telegram-music-bot/venv/bin/python telegram_music_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 